# com.lakeside-motorbikes.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lakeside-motorbikes.worker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/shea/Repos/lakeside-motorbikes/.venv/bin/python</string>
        <string>-m</string>
        <string>lakeside_motorbikes</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/shea/Repos/lakeside-motorbikes</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/shea/Repos/lakeside-motorbikes/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/shea/Repos/lakeside-motorbikes/logs/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
